# systemd/user/background-thumbnails.path
[Unit]
Description=Watch backgrounds for Walker thumbnail cache

[Path]
PathChanged=%h/pictures/backgrounds
PathModified=%h/pictures/backgrounds
Unit=background-thumbnails.service

[Install]
WantedBy=default.target
