# systemd/user/background-thumbnails.service
[Unit]
Description=Cache Walker background thumbnails

[Service]
Type=oneshot
ExecStart=%h/dotfiles/bin/ron-background-generate-thumbnails
